# relay_bot/relay-bot.service
[Unit]
Description=Feishu Relay Bot
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/feishu-relay-bot run --config /etc/relay-bot/config.yaml
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
